# com.jobscraper.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jobscraper.monitor</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/thomassee/Docker/containers/job-scraper/startup-monitor.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>/Users/thomassee/Docker/containers/job-scraper/logs/launchagent.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/thomassee/Docker/containers/job-scraper/logs/launchagent.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/thomassee/Docker/containers/job-scraper</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
    
    <key>StartInterval</key>
    <integer>3600</integer>
    
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>